# schedulers/tap-agent/DEBIAN/postinst
#!/bin/sh

systemctl enable monroe-tap-agent || true
systemctl start  monroe-tap-agent || true
