# schedulers/tap-agent/DEBIAN/postinst
#!/bin/sh
# Generate a new cert since Debian 9 does not work with adhoc mode (too old "werkzeug")
openssl req -x509 -newkey rsa:4096 -nodes -out /opt/monroe/tap-agent/cert.pem -keyout /opt/monroe/tap-agent/key.pem -days 3650 -subj "/O=Monroe Project/CN=localhost"
#Change the default to point to the newly created files
sed -i "s#TAP_AGENT_CERT=.*#TAP_AGENT_CERT=/opt/monroe/tap-agent/cert.pem#g" /opt/monroe/tap-agent/monroe-tap-agent
sed -i "s#TAP_AGENT_KEY=.*#TAP_AGENT_KEY=/opt/monroe/tap-agent/key.pem#g" /opt/monroe/tap-agent/monroe-tap-agent

systemctl enable monroe-tap-agent 
systemctl start  monroe-tap-agent
/usr/sbin/circle restart
