# schedulers/tap-agent/DEBIAN/prerm
#!/bin/sh

systemctl stop monroe-tap-agent 
systemctl disable monroe-tap-agent
