# core/lib/systemd/system/monroe-namespace.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=MONROE Network Namespace
After=network.target docker.service

[Service]
# Consider to use: 
# ExecStartPost=<start up dating monroe-namespace with IF>
ExecStartPre=/usr/bin/pull-base-containers
ExecStart=/usr/bin/monroe-namespace
# Should restrict to only monroe containers
ExecStop=/usr/bin/monroe-namespace stop
Type=simple
Restart=on-failure
TimeoutSec=30min
RestartSec=5
