# core/DEBIAN/postinst
#!/bin/sh
set -e

#To reflect the updated config
systemctl restart docker

#Renews ip on alla interfaces (check so it actuall works on multiple interfaces) https://bit.ly/2G1z7Ib
dhclient -r && dhclient || true
netplan apply || true

mkdir -p /experiments/user
systemctl enable monroe-namespace
systemctl restart monroe-namespace
