# core/DEBIAN/prerm
#!/bin/sh

# TODO: Remove monroe-experiments containers, images and files?

systemctl stop monroe-namespace
systemctl disable monroe-namespace
